# Lab_guardian/debian/postinst
#!/bin/bash
# postinst — runs after the .deb package is unpacked
# Configures auditd rules for terminal command capture and
# installs a .desktop launcher for the application menu.
set -e

GREEN='\033[0;32m'
YELLOW='\033[1;33m'
CYAN='\033[0;36m'
NC='\033[0m'

# ---------------------------------------------------------------------------
# 1. auditd rules for terminal command capture (OPTIONAL)
#    auditd is a Recommended package — the app works without it, but
#    terminal command monitoring (Layer 2) will be disabled.
# ---------------------------------------------------------------------------
RULES_DIR="/etc/audit/rules.d"
RULES_FILE="${RULES_DIR}/lab_guardian_exam.rules"

if command -v auditd > /dev/null 2>&1 || [ -d "$RULES_DIR" ]; then
    mkdir -p "$RULES_DIR"
    cat > "$RULES_FILE" << 'AUDITD_RULES'
## Lab Guardian — monitor network-capable terminal commands
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/curl    -k exam_net
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/wget    -k exam_net
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/git     -k exam_net
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/ssh     -k exam_net
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/python3 -k exam_net
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/pip3    -k exam_net
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/nc      -k exam_net
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/ncat    -k exam_net
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/socat   -k exam_net
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/node    -k exam_net
-a always,exit -F arch=b64 -S execve -F exe=/usr/bin/npm     -k exam_net
AUDITD_RULES
    echo -e "${GREEN}[OK]${NC} auditd rules installed: ${RULES_FILE}"

    # Load rules immediately if auditctl is available
    if command -v auditctl > /dev/null 2>&1; then
        auditctl -R "$RULES_FILE" 2>/dev/null || true
        echo -e "${GREEN}[OK]${NC} auditd rules loaded into kernel"
    fi

    # Ensure auditd is running and enabled
    if command -v systemctl > /dev/null 2>&1 && systemctl is-system-running --quiet 2>/dev/null; then
        systemctl enable auditd 2>/dev/null || true
        systemctl restart auditd 2>/dev/null || true
        echo -e "${GREEN}[OK]${NC} auditd service enabled and restarted"
    fi

    # Make the audit log readable by non-root users so Lab Guardian can tail it
    AUDIT_LOG="/var/log/audit/audit.log"
    if [ -f "$AUDIT_LOG" ]; then
        chmod o+r "$AUDIT_LOG"
        echo -e "${GREEN}[OK]${NC} Audit log made readable: ${AUDIT_LOG}"
    fi
else
    echo -e "${YELLOW}[NOTE]${NC} auditd is not installed — terminal command monitoring (Layer 2) will be disabled."
    echo -e "${YELLOW}[NOTE]${NC} To enable full monitoring later, run:  sudo apt install auditd"
fi

# ---------------------------------------------------------------------------
# 2. Update desktop application database
# ---------------------------------------------------------------------------
if command -v update-desktop-database > /dev/null 2>&1; then
    update-desktop-database /usr/share/applications 2>/dev/null || true
fi

echo ""
echo -e "${GREEN}╔══════════════════════════════════════════════════════╗${NC}"
echo -e "${GREEN}║   Lab Guardian installed successfully!               ║${NC}"
echo -e "${GREEN}║                                                      ║${NC}"
echo -e "${GREEN}║   Launch: lab_guardian   (terminal)                  ║${NC}"
echo -e "${GREEN}║       or: search 'LabGuardian' in the app menu       ║${NC}"
echo -e "${GREEN}╚══════════════════════════════════════════════════════╝${NC}"
echo ""
